# scripts/system/synology-health-check.service
[Unit]
Description=Synology NFS Health-Check
Before=borgbackup.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/synology-health-check.sh
